# root/etc/s6-overlay/s6-rc.d/init-mod-calibre-worddumb-pip-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-calibre-worddumb-pip-add-package/run